# service/systemd/conet-prysm-validator.service
[Unit]
Description=CoNET Prysm validator client (node-0)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=peter
Group=peter
WorkingDirectory=/home/peter/ethereum-pos-mainnet
EnvironmentFile=-/etc/default/conet-validator-redeem-listener
ExecStart=/home/peter/ethereum-pos-mainnet/10_run_prysm_validator.sh
Restart=always
RestartSec=10
TimeoutStopSec=120
StandardOutput=append:/home/peter/ethereum-pos-mainnet/network/node-0/logs/validator.log
StandardError=append:/home/peter/ethereum-pos-mainnet/network/node-0/logs/validator.log
SyslogIdentifier=conet-prysm-validator
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
