# service/systemd/ipfs-conet.base-op-reth-native.service
[Unit]
Description=Base native op-reth (ipfs.conet.network)
After=network-online.target
Wants=network-online.target

[Service]
User=peter
Group=peter
WorkingDirectory=/home/peter
Environment=RUST_LOG=info
# NUMA: bind CPUs and memory to nodes 0+1 (nodes 2–3 have no local memory on this host)
# Static files: --datadir.static-files. ETL temp: reth.toml [stages.etl] dir = /media/sda/baseTemp (not a CLI flag).
ExecStart=/usr/bin/numactl --cpunodebind=0,1 --membind=0,1 -- /media/sda/bin/op-reth-v1.11.3 node --datadir=/home/peter/base/op-reth --datadir.static-files=/media/sda/base/op-reth-static --chain=base --rollup.sequencer-http=https://mainnet-sequencer.base.org --rollup.disable-tx-pool-gossip --http --http.addr=0.0.0.0 --http.port=8547 --http.api=eth,net,web3,txpool --http.corsdomain=* --ws --ws.addr=0.0.0.0 --ws.port=8548 --ws.origins=* --ws.api=web3,eth,net,txpool --authrpc.addr=0.0.0.0 --authrpc.port=8552 --authrpc.jwtsecret=/home/peter/base/jwt/jwt.hex --port=30313 --discovery.port=30313 --addr=0.0.0.0 --discovery.addr=0.0.0.0 --max-outbound-peers=100 --nat=extip:38.102.126.33 --bootnodes=enode://87a32fd13bd596b2ffca97020e31aef4ddcc1bbd4b95bb633d16c1329f654f34049ed240a36b449fda5e5225d70fe40bc667f53c304b71f8e68fc9d448690b51@3.231.138.188:30301,enode://ca21ea8f176adb2e229ce2d700830c844af0ea941a1d8152a9513b966fe525e809c3a6c73a2c18a12b74ed6ec4380edf91662778fe0b79f6a591236e49e176f9@184.72.129.189:30301,enode://acf4507a211ba7c1e52cdf4eef62cdc3c32e7c9c47998954f7ba024026f9a6b2150cd3f0b734d9c78e507ab70d59ba61dfe5c45e1078c7ad0775fb251d7735a2@3.220.145.177:30301,enode://8a5a5006159bf079d06a04e5eceab2a1ce6e0f721875b2a9c96905336219dbe14203d38f70f3754686a6324f786c2f9852d8c0dd3adac2d080f4db35efc678c5@3.231.11.52:30301,enode://cdadbe835308ad3557f9a1de8db411da1a260a98f8421d62da90e71da66e55e98aaa8e90aa7ce01b408a54e4bd2253d701218081ded3dbe5efbbc7b41d7cef79@54.198.153.150:30301 --db.max-readers=32767 --engine.cross-block-cache-size=32768 --engine.memory-block-buffer-target=4 --engine.parallel-sparse-trie --rpc.max-connections=600 --rpc.gascap=18446744073709551615
Restart=always
RestartSec=5
LimitNOFILE=1048576

[Install]
WantedBy=multi-user.target
